# config.example.yaml
splunk:
  host: localhost
  rest_port: 8089        # REST API / management port
  hec_port: 8088         # HTTP Event Collector port
  hec_token: "YOUR_HEC_TOKEN_HERE"
  username: admin        # local Splunk Enterprise admin (not your splunk.com email)
  password: "YOUR_SPLUNK_PASSWORD_HERE"
  index_baseline: arena_baseline
  index_attacks: arena_attacks
  verify_ssl: false      # self-signed cert on local installs
  # MCP token — generate via setup_verify.py after installing the MCP Server app.
  # Leave blank if you haven't installed the MCP app yet (system still works via REST).
  mcp_token: ""

llm:
  # Options: gemini | ollama | splunk_hosted
  # gemini  — free API key from aistudio.google.com, fastest to get started
  # ollama  — fully local, requires `ollama serve` + model pull
  # splunk_hosted — Foundation-sec-1.1 via Splunk Cloud, earns Hosted Models prize
  provider: gemini

  gemini:
    api_key: "YOUR_GEMINI_API_KEY"
    model: gemini-2.0-flash

  ollama:
    model: qwen2.5:3b    # faster, lower RAM; run: ollama pull qwen2.5:3b
    # model: llama3.1    # larger, slower; run: ollama pull llama3.1
    base_url: http://localhost:11434

  splunk_hosted:
    endpoint: "https://YOUR_INSTANCE.splunkcloud.com/services/ml/models/foundation-sec-1.1-8b-instruct/predict"
    token: "YOUR_SPLUNK_CLOUD_TOKEN"

  groq:
    # Free tier: 6000 req/day, 30 req/min — get key at console.groq.com (no credit card)
    api_key: "YOUR_GROQ_API_KEY"
    model: llama-3.3-70b-versatile   # best quality; swap to llama-3.1-8b-instant for speed

  # Track 4 — LLM data security. When present, every prompt is PII-stripped
  # (usernames, internal IPs, hostnames, domains, GUIDs, SIDs, user paths) before
  # it leaves the process and rehydrated in the rule that comes back; a hash-only,
  # tamper-evident audit record is written per call. Omit/false to disable.
  # For true on-prem privacy use provider: ollama (nothing leaves the network).
  secure:
    audit: true
    audit_path: results/llm_audit.jsonl
    store_preview: false       # true only in dev — keeps a 200-char preview, not hash-only

arena:
  num_rounds: 5
  indexing_wait_seconds: 4   # seconds to wait after HEC inject before searching
  techniques:
    - T1110.001   # Brute Force: Password Guessing
    - T1021.001   # Remote Services: RDP
    - T1053.005   # Scheduled Task / Job
    - T1136.001   # Create Account: Local Account
    - T1003.001   # OS Credential Dumping: LSASS
    - T1547.001   # Boot/Logon Autostart: Registry Run Keys
    - T1059.001   # PowerShell: Encoded Command
    - T1055.001   # Process Injection: DLL Injection
    - T1562.001   # Impair Defenses: Disable Windows Defender
    - T1078.004   # Valid Accounts: Cloud Accounts (Azure AD / Entra ID) — Red Canary 2025 #1
    - T1114.003   # Email Collection: Email Forwarding Rule (M365 / BEC)

  # Win conditions (turn-based and real-time modes)
  win_conditions:
    blue_coverage_pct: 70.0    # Blue wins by holding ≥70% coverage for the full duration
    red_compromise_pct: 60.0   # Red wins by compromising (evade + burn) ≥60% of techniques

  # Real-time engine tuning
  realtime:
    red_base_seconds: 6.0        # Red's normal attack cadence
    red_urgent_seconds: 2.5      # Red's cadence after Blue catches a technique
    blue_base_seconds: 4.0       # Blue's normal detection sweep cadence
    blue_urgent_seconds: 2.0     # Blue's cadence after a miss
    window_seconds: 120          # sliding detection window width
    evasion_grace_seconds: 12.0  # seconds before a missed technique counts as evaded
    metrics_every_seconds: 8.0   # scoreboard + dashboard feed interval
    benign_every_seconds: 12.0   # benign traffic injection interval
    poison_interval_seconds: 15.0 # Red's FP-flooding cadence (alert fatigue)
    infer_mutations: false       # Track 5 moat: when true, Blue INFERS Red's field
                                 # changes by diffing evading events vs. the last batch
                                 # it caught — instead of being handed Red's overrides
                                 # (production-realistic; the attacker never confesses)

# ──────────────────────────────────────────────────────────────────────────────
# PRODUCTION INTEGRATIONS (all OPTIONAL, all default-OFF)
# The arena above is fully self-contained. These stanzas wire the production-hardening
# modules to real infrastructure once you provision it. Each requires external setup —
# see the run-book at the bottom of the named source file. Leaving enabled:false keeps
# the demo byte-for-byte unchanged.
# ──────────────────────────────────────────────────────────────────────────────

# Track 1 — real attack data via MITRE Caldera (red_agent/caldera_client.py).
# Replaces synthetic HEC injection with real adversary emulation on lab hosts whose
# Sysmon/forwarder telemetry flows into index_attacks the normal way.
caldera:
  enabled: false
  base_url: "http://YOUR_CALDERA_HOST:8888"
  api_key: "YOUR_CALDERA_RED_API_KEY"     # conf/local.yml -> api_key_red
  adversary_id: "YOUR_ADVERSARY_PROFILE_ID"
  group: ""                                # "" = all agents
  verify_ssl: false

# Track 5 L3 — EDR ground truth (blue_agent/edr_client.py). A second, authoritative
# witness to what executed on the endpoint; corroborate() turns its view into Blue's
# blind-spot list. CrowdStrike Falcon schema (Defender works with schema tweaks).
edr:
  enabled: false
  base_url: "https://api.crowdstrike.com"  # or your Falcon cloud, e.g. api.us-2.crowdstrike.com
  client_id: "YOUR_FALCON_CLIENT_ID"
  client_secret: "YOUR_FALCON_CLIENT_SECRET"

# Track 5 L4 — proactive detections from threat intel (blue_agent/threat_intel.py).
# Point at advisory text/URLs or ATT&CK Navigator layer files; gap_analysis ranks the
# techniques Blue should pre-build rules for. Offline-capable; no credentials needed.
threat_intel:
  enabled: false
  navigator_layers: []     # e.g. ["intel/apt29_layer.json"]
  advisory_urls: []        # e.g. ["https://www.cisa.gov/.../aa24-xxxa"]
