# red_agent/templates/T1053.005.json
{
  "technique_id": "T1053.005",
  "technique_name": "Scheduled Task/Job: Scheduled Task",
  "tactic": "Persistence",
  "events": [
    {
      "note": "schtasks.exe process creation — LOLBin used to register the task",
      "sourcetype": "XmlWinEventLog:Microsoft-Windows-Sysmon/Operational",
      "template": {
        "EventCode": 1,
        "EventDescription": "Process Create",
        "Image": "C:\\Windows\\System32\\schtasks.exe",
        "CommandLine": "schtasks.exe /Create /TN \"{task_name}\" /TR \"{malicious_command}\" /SC ONLOGON /RU SYSTEM /F",
        "ParentImage": "C:\\Windows\\System32\\cmd.exe",
        "User": "CORP\\{creating_user}",
        "IntegrityLevel": "High",
        "Hashes": "SHA256=76D79B2B73C6B72B4B0B8B3B8B3B8B3B8B3B8B3B8B3B8B3B8B3B8B3B8B3B8B3"
      },
      "count": 1,
      "spread_seconds": 1
    },
    {
      "note": "Windows Security 4698 — task creation audit event with full XML TaskContent",
      "sourcetype": "WinEventLog:Security",
      "template": {
        "EventCode": 4698,
        "Account_Name": "{creating_user}",
        "Account_Domain": "CORP",
        "TaskName": "{task_name}",
        "TaskContent": "<Task><Triggers><LogonTrigger><Enabled>true</Enabled></LogonTrigger></Triggers><Actions><Exec><Command>{malicious_command}</Command></Exec></Actions><Principals><Principal><RunLevel>HighestAvailable</RunLevel></Principal></Principals></Task>"
      },
      "count": 1,
      "spread_seconds": 2
    },
    {
      "note": "Payload executes under svchost.exe — confirms task ran",
      "sourcetype": "XmlWinEventLog:Microsoft-Windows-Sysmon/Operational",
      "template": {
        "EventCode": 1,
        "EventDescription": "Process Create",
        "Image": "C:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe",
        "CommandLine": "powershell.exe -NonI -W Hidden -Exec Bypass -Enc JABjAD0ATgBlAHcALQBPAGIAagBlAGMAdAAoAA==",
        "ParentImage": "C:\\Windows\\System32\\svchost.exe",
        "User": "NT AUTHORITY\\SYSTEM",
        "IntegrityLevel": "System"
      },
      "count": 1,
      "spread_seconds": 5
    }
  ],
  "mutation_hints": {
    "evadable_fields": [
      "TaskName",
      "CommandLine",
      "ParentImage",
      "IntegrityLevel",
      "Image",
      "TaskContent",
      "User",
      "Hashes",
      "Account_Domain"
    ],
    "objective": "Maintain SYSTEM-level persistence that survives user logoff and reboots",
    "dwell_threshold_seconds": 120,
    "detection_anchors": [
      "Field match: Image=schtasks.exe → use C:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe with Register-ScheduledTask cmdlet; no schtasks.exe process is spawned, evading schtasks-specific rules",
      "Field match: CommandLine contains /Create or schtasks → PowerShell syntax has no /Create flag; also at.exe (legacy) uses 'at HH:MM cmd' syntax",
      "Field match: IntegrityLevel=High → reduce to Medium if UAC was already bypassed upstream; High-integrity task creation is a key alert",
      "Field match: ParentImage=cmd.exe → change to explorer.exe, powershell.exe, or wscript.exe as the parent spawning the task registration",
      "XML TaskContent trigger: LogonTrigger → use BootTrigger (on system start) or a TimeTrigger with repetition interval — less alarming trigger types",
      "Field match: TaskContent contains HighestAvailable RunLevel → use LeastPrivilege in the Principal to appear benign",
      "Field: User → change from CORP\\creating_user to NT AUTHORITY\\SYSTEM if already elevated, or use a service account name"
    ]
  }
}
