# red_agent/templates/T1547.001.json
{
  "technique_id": "T1547.001",
  "technique_name": "Boot/Logon Autostart: Registry Run Keys",
  "tactic": "Persistence",
  "events": [
    {
      "note": "reg.exe used to write the run key — catches both API and CLI-based persistence",
      "sourcetype": "XmlWinEventLog:Microsoft-Windows-Sysmon/Operational",
      "template": {
        "EventCode": 1,
        "EventDescription": "Process Create",
        "Image": "C:\\Windows\\System32\\reg.exe",
        "CommandLine": "reg add \"HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\" /v {reg_key_name} /t REG_SZ /d \"{malicious_command}\" /f",
        "ParentImage": "{setting_process}",
        "User": "CORP\\{target_user}",
        "IntegrityLevel": "Medium"
      },
      "count": 1,
      "spread_seconds": 1
    },
    {
      "note": "Sysmon 13: registry value set — catches API-based writes that don't use reg.exe",
      "sourcetype": "XmlWinEventLog:Microsoft-Windows-Sysmon/Operational",
      "template": {
        "EventCode": 13,
        "EventDescription": "Registry value set",
        "EventType": "SetValue",
        "Image": "{setting_process}",
        "TargetObject": "HKU\\S-1-5-21-1234567890-1234567890-1234567890-1001\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\{reg_key_name}",
        "Details": "{malicious_command}",
        "User": "CORP\\{target_user}"
      },
      "count": 1,
      "spread_seconds": 1
    }
  ],
  "mutation_hints": {
    "evadable_fields": [
      "TargetObject",
      "Details",
      "Image",
      "CommandLine",
      "ParentImage",
      "User",
      "IntegrityLevel"
    ],
    "objective": "Execute malicious payload on every user logon without UAC prompt",
    "dwell_threshold_seconds": 120,
    "detection_anchors": [
      "Field match: Image=reg.exe → use C:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe with Set-ItemProperty cmdlet; direct .NET registry API calls don't spawn reg.exe at all",
      "Field match: TargetObject contains \\CurrentVersion\\Run\\ → use \\CurrentVersion\\RunOnce\\ (runs once then deletes itself, less persistent but avoids Run-key rules), or HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Userinit (different autorun path)",
      "Wildcard match: TargetObject contains HKCU or HKLM followed by Run → use HKCU\\Environment\\UserInitMprLogonScript (logon script via environment variable; less commonly monitored)",
      "Field: Details (the payload value) → use a signed binary with malicious arguments, or a path that looks like a legitimate update binary",
      "Field: ParentImage → change the process that invokes reg.exe; cmd.exe is most suspicious; explorer.exe or svchost.exe as parent looks more benign",
      "Field: IntegrityLevel=Medium → already at Medium for HKCU writes; High would be needed for HKLM; if switching to HKLM path, raise to High and adjust TargetObject accordingly"
    ]
  }
}
