# blue_agent/rules/baseline/scheduled_task.spl
index=arena_attacks (EventCode=4698) OR (EventCode=1 Image="*\\schtasks.exe" CommandLine="*/Create*") OR (EventCode=1 IntegrityLevel="System" (ParentImage="*\\taskeng.exe" OR ParentImage="*\\svchost.exe") CommandLine="*powershell*")
| eval technique="T1053.005", rule_name="scheduled_task_baseline"
