# blue_agent/rules/baseline/scheduled_task.spl
index=arena_attacks EventCode=1 Image="*\\schtasks.exe" CommandLine="*/Create*"
| eval technique="T1053.005", rule_name="scheduled_task_baseline"
